# crontab/racked-pc.crontab
# ##############################################################################
# pc.crontab
#
# This file is used to backup my racked pc to my NAS.
# ##############################################################################

00 10 * * * ssh iolave@127.0.0.1 "/opt/iolave/crontab-backup.sh"
